# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/var/log/supervisord.log ; log file; default $CWD/supervisord.log
logfile_maxbytes=50MB ; log file max bytes b4 rotation; default 50MB
logfile_backups=10 ; number of log file backups; default 10
loglevel=info ; log level; default info; other levels: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
; user=chrisk ; (default is no user)
nodaemon=false ; start in foreground if true; default false

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx_error.log
stdout_logfile=/var/log/nginx_access.log

[program:php-fpm]
command=/usr/sbin/php-fpm7.3 --nodaemonize
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm_error.log
stdout_logfile=/var/log/php-fpm_access.log
